# data/cron.d/c2-gs-aggregate-billing
# Run every minute
*/1 * * * * root c2-gs-run aggregate-billing
